# extra/10-local.rules
# Dualshock 4 LED
ACTION=="add", SUBSYSTEM=="input", ATTRS{name}=="Wireless Controller" RUN+="/opt/ds4perm '%p'"
